# files/mx_master2_config/xbindkeys.service
[Unit]
Description=Start xbindkeys

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/bin/xbindkeys

[Install]
WantedBy=multi-user.target
